# crontab
57 * * * *	/usr/local/bin/firetv-activate-remap
